# muslsys/nginx/src/main/dist/nginx.service
[Unit]
Description=nginx
Requires=network.target
After=syslog.target

[Service]
Type=forking
PIDFile=/var/log/nginx/nginx.pid
ExecStartPre=/usr/libexec/nginx/nginx -c /etc/nginx/nginx.conf -t
ExecStart=/usr/libexec/nginx/nginx -c /etc/nginx/nginx.conf
ExecReload=/usr/libexec/nginx/nginx -s reload
ExecStop=/bin/kill -s QUIT $MAINPID
PrivateTmp=true
User=root
Group=root

[Install]
WantedBy=default.target
